# templates/gem_server_initd.sh.erb
#! /bin/sh
### BEGIN INIT INFO
# Provides:               gem
# Required-Start:         $network
# Required-Stop:
# Default-Start:          <%= @gem_run_levels.class == Array ? @gem_run_levels.join(' ') : @gem_run_levels %>
# Default-Stop:           <%= @gem_not_run_levels.class == Array ? @gem_not_run_levels.join(' ') : @gem_not_run_levels %>
# Short-Description:      <%= @gem_server_description %>
# Description:            <%= @gem_server_description %>
### END INIT INFO

# Check for missing binaries
GEM_BIN=<% =@gem_path %>
test -x $GEM_BIN || { echo "$GEM_BIN not installed";
        if [ "$1" = "stop" ]; then exit 0;
        else exit 5; fi; }
